# systemd/xray-update.service
[Unit]
Description=Update Xray and geodata

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-xray.sh
